# HelloWorld4/package/DEBIAN/prerm
#!/bin/sh
set -e

systemctl stop helloworldapp.service || true
systemctl disable helloworldapp.service || true

exit 0
